# init/systemd/extractarr.service
[Unit]
Description=Extractarr Service

[Service]
ExecStart=/usr/local/bin/extractarr
Restart=always

[Install]
WantedBy=multi-user.target